# contrib/user/vec-updatedb.service
[Unit]
Description=vec semantic index update (userland)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/vec --config %h/.config/vec/config.toml updatedb
ExecStartPost=%h/.local/bin/vec --config %h/.config/vec/config.toml gc
Nice=19
IOSchedulingClass=idle
# Allow enough time for full home directory indexing
TimeoutStartSec=1h

[Install]
WantedBy=default.target
